# sleepcare-core/files/sleepcare-core.service
[Unit]
Description=SleepCare Core Daemon (sleepcare-core)
Documentation=https://github.com/WhiPaper/sleepcare
After=network-online.target earsys.service avahi-daemon.service
Wants=network-online.target avahi-daemon.service
Requires=earsys.service

[Service]
Type=simple
RuntimeDirectory=sleepcare
RuntimeDirectoryMode=0755
ConfigurationDirectory=sleepcare
ReadWritePaths=/etc/sleepcare
ExecStart=/usr/bin/sleepcare-core
Restart=on-failure
RestartSec=3s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sleepcare-core

# Security hardening
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
